# ai-frontend.service
[Unit]
Description=AI Frontend Service
After=network.target

[Service]
Type=simple
User=yjw
WorkingDirectory=/home/yjw/Ai-1.0
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10
Environment=NODE_ENV=development

[Install]
WantedBy=multi-user.target
